# cookbooks/web_app/templates/default/webapp.conf.erb
description 'nodejs upstart script'

start on filesystem and started networking
respawn

chdir /web_app

exec npm start